# cluster/partition.service
[Unit]
Description=Beam Partition

[Service]
User=stack
Restart=on-failure
WorkingDirectory=/home/stack/protobeam/
ExecStart=/home/stack/protobeam/cluster/startpartition.sh -cfg c3.json

[Install]
WantedBy=multi-user.target
